# system_files/usr/share/containers/systemd/jellyfin.container
[Container]
ContainerName=jellyfin
Image=docker.io/jellyfin/jellyfin:latest
AutoUpdate=registry
PublishPort=8096:8096/tcp
UserNS=keep-id

# Internal volumes (managed by Podman)
#Volume=jellyfin-config:/config:Z
#Volume=jellyfin-cache:/cache:Z
Volume=/mnt/orico_a/containers/jellyfin/config:/config:Z
Volume=/mnt/orico_a/containers/jellyfin/cache:/cache:Z

# External bind mount for media
Volume=/mnt/orico_a/media_movies:/movies:ro,Z
Volume=/mnt/orico_a/media_shows:/shows:ro,Z
Volume=/mnt/orico_a/media_videos:/videos:ro,Z

[Service]
# Make systemd treat exit code 143 (SIGTERM) as normal
SuccessExitStatus=0 143

[Install]
#WantedBy=multi-user.target
WantedBy=default.target
